# 9.3/cron.d/backup-database
0 2 * * * gosu postgres /usr/bin/envdir /etc/wal-e.d/env /usr/local/bin/wal-e backup-push ${PGDATA}
