# 9.3/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/local/bin/run-cron.sh
autostart=true
autorestart=true
user=root
startsecs=3
