# 9.3/supervisor/conf.d/postgres.conf
[program:postgres]
command=/usr/local/bin/run-postgres.sh postgres
autostart=true
autorestart=true
user=root
startsecs=5
